# cf-ddns.timer
[Unit]
Description=Runs cf-ddns script every minute

[Timer]
OnCalendar=*:0/1

[Install]
WantedBy=timers.target
